# com.tuatha.autopush.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tuatha.autopush</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/user/Desktop/tuatha-patient-app/auto-push.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>StandardOutPath</key>
    <string>/Users/user/Desktop/tuatha-patient-app/autopush.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/user/Desktop/tuatha-patient-app/autopush-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
